# systemd/user/mako.service
[Unit]
	Description=Mako Notify Daemon

[Service]
	ExecStart=/usr/bin/mako 

[Install]
	WantedBy=default.target
