# systemd/user/mpd.service
[Unit]
	Description=MPD Daemon

[Service]
	ExecStart=/usr/bin/mpd --no-daemon

[Install]
	WantedBy=default.target
